# rc.local
#!/bin/bash
cd "$(dirname "$0")"
cp rpc/* /etc/qubes-rpc/
